# recipes-extended/voice-assistant-switcher/voice-assistant-switcher/vad.timer
[Unit]
Description=Delay the start of vad at start-up

[Timer]
OnBootSec=12s

[Install]
WantedBy=timers.target

